# packages/sys-devel/distcc/files/distcc-initscript
#!/sbin/runscript
# Copyright 2009 Jan Meier
# Copyright 1999-2009 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	need net
}

start() {
	ebegin "Starting distccd"

	if [ ! -e /var/run/distccd ] ; then
		mkdir -p /var/run/distccd
		chown distcc:daemon /var/run/distccd
	fi

	start-stop-daemon --start --quiet --exec "/usr/bin/distccd" -- \
		--daemon --pid-file "/var/run/distccd/pid" --user distcc \
		${DISTCCD_OPTS}
	eend $?
}

stop() {
	ebegin "Stopping distccd"
	start-stop-daemon --stop --quiet --pidfile "/var/run/distccd/pid"
	eend $?
}
